# docker/etc/supervisord.d/processes.ini
[program:dockergen]
command=/usr/bin/docker-gen -notify "/usr/bin/supervisorctl restart filebeat" -notify-output -watch /etc/docker-gen/filebeat.tmpl /etc/filebeat/filebeat.yml
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:filebeat]
command=/usr/share/filebeat/bin/filebeat -c /etc/filebeat/filebeat.yml -path.home /usr/share/filebeat -path.config /etc/filebeat -path.data /var/lib/filebeat -path.logs /var/log/filebeat
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
